# systemd/z13-hibernate-gate.service
[Unit]
Description=z13 hibernate gate pre-check (quiesce and LED state before systemd-hibernate)
Before=systemd-hibernate.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/lib/z13-hibernate/gate-hook.sh
RemainAfterExit=yes
TimeoutStartSec=180

[Install]
RequiredBy=hibernate.target
