# systemd/systemd-hibernate.service.d/10-gate.conf
[Unit]
# Ensures the gate pre-check runs and completes before we proceed with hibernate.
Requires=z13-hibernate-gate.service
After=z13-hibernate-gate.service
